# root/etc/s6-overlay/s6-rc.d/probot/run
#!/command/with-contenv bash
cd /app

exec \
	s6-setuidgid abc npm --prefix /app --cache /app start